# mqtt_daemon/r4dcb08-mqtt.service
[Unit]
Description=R4DCB08 Temperature Sensor MQTT Publisher
After=network-online.target mosquitto.service
Wants=network-online.target mosquitto.service
StartLimitIntervalSec=300
StartLimitBurst=5
ConditionPathExists=/etc/r4dcb08-mqtt.conf

[Service]
Type=notify
NotifyAccess=main
ExecStart=/usr/local/bin/r4dcb08-mqtt -c /etc/r4dcb08-mqtt.conf
ExecReload=/bin/kill -HUP $MAINPID
TimeoutStartSec=120
WatchdogSec=3600
Restart=on-failure
RestartSec=5

# Run as dedicated user (create with: useradd -r -s /usr/sbin/nologin r4dcb08)
User=r4dcb08
Group=r4dcb08
SupplementaryGroups=dialout

# Security hardening
ReadWritePaths=/var/run
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
PrivateDevices=false

# Allow access to serial devices
DeviceAllow=/dev/ttyUSB* rw
DeviceAllow=/dev/ttyACM* rw
DeviceAllow=/dev/ttyAMA* rw
DeviceAllow=/dev/serial* rw

[Install]
WantedBy=multi-user.target
